# .config/systemd/user/dailybackup.service
[Unit]
Description=Daily backup service

[Service]
Type=simple
ExecStart=/home/greg/bin/backmeup
